# overlay-rksdk/files-overlay-rk3288/etc/init/ttyFIQ0.conf
start on starting boot-services
stop on stopping boot-services

respawn
exec /sbin/agetty 115200 ttyFIQ0 linux -a root
